# templates/com.auto-transcribe.pull-plaud.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.auto-transcribe.pull-plaud</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__BASE_DIR__/pull-plaud.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__BASE_DIR__</string>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/tmp/plaud-pull-out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/plaud-pull-err.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
